# redroid-install/redroid-webportal-supervisor.conf
[program:redroid-webportal]
directory=/home/redroid/redroid-webportal
command=rails s
autostart=true
autorestart=true
